# b.sh
#!/bin/sh /etc/rc.common

start(){

sh /root/d.sh

}

reload(){

sh /root/g.sh

}
